# install/etc/services.d/webhook/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping webhook"
killall -9 webhook
